# zjenkins/xjm-test/bin/jenkins.service
[Unit]
Description=Jenkins
After=network.target

[Service]
Type=forking
PIDFILE=/tmp/jenkins.pid
ExecStart=/usr/local/bin/jenkins.sh start
ExecReload=
ExecStop=/usr/local/bin/jenkins.sh stop
PrivateTmp=true

[Install]
WantedBy=multi-user.target